# supervisor.ini
[supervisord]
nodaemon = true
user = root

[program: default-app]
command = /bin/bash -c "exec /run_app.sh"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[program: envoy-proxy-template]
user = root
command = /bin/bash -c "exec /run_envoy.sh"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

